# scripts/octopus.service
[Unit]
Description=Koupper Octopus Engine Daemon
After=network.target

[Service]
Type=simple
User=pi
Group=pi
# Path to your installation (adjust if ~/.koupper/bin is elsewhere)
WorkingDirectory=/home/pi
Environment="KOUPPER_OCTOPUS_TOKEN=your-secret-token-here"
Environment="KOUPPER_OCTOPUS_HOST=0.0.0.0"
ExecStart=/home/pi/.koupper/bin/octopus
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
